# infrastructure/crontab
# Rinawarp Platform Maintenance Schedule

# Weekly maintenance - Every Monday at 00:00
0 0 * * 1 cd /opt/rinawarp && ./scripts/maintenance.sh weekly >> /var/log/rinawarp-maintenance/weekly.log 2>&1

# Monthly maintenance - 1st of every month at 01:00
0 1 1 * * cd /opt/rinawarp && ./scripts/maintenance.sh monthly >> /var/log/rinawarp-maintenance/monthly.log 2>&1

# Quarterly maintenance - 1st of Jan/Apr/Jul/Oct at 02:00
0 2 1 1,4,7,10 * cd /opt/rinawarp && ./scripts/maintenance.sh quarterly >> /var/log/rinawarp-maintenance/quarterly.log 2>&1

# Daily backup verification at 03:00
0 3 * * * cd /opt/rinawarp && ./scripts/verify-backups.sh >> /var/log/rinawarp-maintenance/backup-verify.log 2>&1

# Clean old logs weekly - Every Sunday at 04:00
0 4 * * 0 find /var/log/rinawarp-maintenance -type f -name "*.log" -mtime +30 -delete

# Security scan - Every day at 05:00
0 5 * * * cd /opt/rinawarp && ./scripts/security-scan.sh >> /var/log/rinawarp-maintenance/security-scan.log 2>&1

# Test monitoring - Weekly on Sunday at 06:00
0 6 * * 0 cd /opt/rinawarp && ./scripts/test-monitor.js >> /var/log/rinawarp-maintenance/test-monitor.log 2>&1

# Documentation check - Daily at 07:00
0 7 * * * cd /opt/rinawarp && ./scripts/doc-monitor.js >> /var/log/rinawarp-maintenance/doc-monitor.log 2>&1

# Performance monitoring - Every hour
0 * * * * cd /opt/rinawarp && ./scripts/performance-monitor.js >> /var/log/rinawarp-maintenance/performance.log 2>&1

# Error rate monitoring - Every 5 minutes
*/5 * * * * cd /opt/rinawarp && ./scripts/error-monitor.js >> /var/log/rinawarp-maintenance/error.log 2>&1

# Usage pattern analysis - Daily at 08:00
0 8 * * * cd /opt/rinawarp && ./scripts/analyze-usage.js >> /var/log/rinawarp-maintenance/usage.log 2>&1
